# packaging/systemd/field-discovery-backup.service
[Unit]
Description=CodexNet verified SQLite backup
Documentation=file:/opt/field-discovery/docs/storage-recovery.md
After=field-discovery-recovery.service
Wants=field-discovery-recovery.service
StartLimitIntervalSec=1h
StartLimitBurst=3

[Service]
Type=oneshot
User=field-discovery
Group=field-discovery
UMask=0077
ExecStart=/opt/field-discovery/venv/bin/field-discovery --config /etc/field-discovery/config.yaml db backup
TimeoutStartSec=10min
Nice=15
CPUWeight=10
IOWeight=10
ProtectSystem=strict
ProtectHome=true
ReadOnlyPaths=/etc/field-discovery
ReadWritePaths=/var/lib/field-discovery
PrivateTmp=true
PrivateDevices=true
NoNewPrivileges=true
CapabilityBoundingSet=
AmbientCapabilities=
RestrictSUIDSGID=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectClock=true
ProtectHostname=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictRealtime=true
RestrictNamespaces=true
RestrictAddressFamilies=AF_UNIX
SystemCallArchitectures=native
MemoryMax=192M
TasksMax=16
